# roles/otelcol/templates/otelcol-svc.j2
#!/sbin/openrc-run
description="OpenTelemetry Collector"
start() {
  nohup otelcol --config=/root/otelcol.yml &
}
stop() {
  pkill -f otelcol
}